# panel/systemd/realm-panel.service
[Unit]
Description=Realm Pro Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/realm-panel/env
WorkingDirectory=/opt/realm-panel/panel/app
ExecStart=/opt/realm-panel/venv/bin/uvicorn main:app --host 0.0.0.0 --port ${PANEL_PORT}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
