# agent/systemd/realm-agent.service
[Unit]
Description=Realm Pro Agent API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/realm-agent/env
WorkingDirectory=/opt/realm-agent/agent/app
ExecStart=/opt/realm-agent/venv/bin/uvicorn main:app --host 0.0.0.0 --port ${AGENT_PORT}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
